# packaging/app-two-factor-auth.cron
*/5 * * * * root /usr/clearos/apps/two_factor_auth/deploy/recycle_tokens > /dev/null 2>&1
